# ansible/curl_job.timer.j2
[Unit]
Description=Timer for {{ service_description }}

[Timer]
OnCalendar={{ timer_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
