# ansible/curl_job.service.j2
[Unit]
Description={{ service_description }}

[Service]
Type=oneshot
ExecStart={{ service_exec }}
